# backend/infra/srv/service/srv.service
[Unit]
Description=Zyg Srv Go Application
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
ExecStart=/usr/local/bin/app
WorkingDirectory=/usr/local/bin
Environment=ENV=production
Restart=always
RestartSec=5

NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

# Additional systemd features for better service management
MemoryAccounting=true
CPUAccounting=true
LimitNOFILE=65536
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target

